# profile.d/fzf.sh
#!/bin/sh
[ -f ~/.fzf.zsh ] && source ~/.fzf.zsh
