# profile.d/travis.sh
#!/bin/sh
[ -f /home/odedlaz/.travis/travis.sh ] && source /home/odedlaz/.travis/travis.sh
